# pkg/DEBIAN/postinst
#!/bin/bash
set -e

# Permissions for evdev / uinput
modprobe uinput || true
echo uinput > /etc/modules-load.d/uinput.conf
echo 'KERNEL=="uinput", MODE="0660", GROUP="input", OPTIONS+="static_node=uinput"' > /etc/udev/rules.d/99-uinput.rules
udevadm control --reload-rules || true
udevadm trigger || true

BACKEND_DIR=/opt/linuxremoteplayer/backend
cd "$BACKEND_DIR"
if [ ! -d ".venv" ]; then
    python3 -m venv .venv
fi
set +e
if ! .venv/bin/pip install -r requirements.txt; then
    echo "[!] Advertencia: falló pip install. Dependencias incompletas."
    touch /opt/linuxremoteplayer/.deps_incomplete
fi
set -e

# Create lrp-setup
cat <<'INNEREOF' > /usr/local/bin/lrp-setup
#!/bin/bash
if [ "$EUID" -ne 0 ]; then
  echo "Please run as root (sudo lrp-setup)"
  exit 1
fi
cd /opt/linuxremoteplayer/scripts
./install.sh
INNEREOF
chmod 755 /usr/local/bin/lrp-setup

# Create lrp-update wrapper
cat <<'INNEREOF' > /usr/local/bin/lrp-update
#!/bin/bash
set -e
export DEBIAN_FRONTEND=noninteractive
MANIFEST_URL="${UPDATE_MANIFEST_URL:-https://linux-remote-player.vercel.app/latest.json}"
echo "Checking $MANIFEST_URL for updates..."
LATEST_JSON=$(curl -fsSL "$MANIFEST_URL")
VERSION=$(echo "$LATEST_JSON" | grep -o '"version"\s*:\s*"[^"]*"' | cut -d'"' -f4)
DEB_URL=$(echo "$LATEST_JSON" | grep -o '"deb_url"\s*:\s*"[^"]*"' | cut -d'"' -f4)
SHA256=$(echo "$LATEST_JSON" | grep -o '"sha256"\s*:\s*"[^"]*"' | cut -d'"' -f4)

cd /tmp
echo "Downloading $DEB_URL..."
curl -fsSL -o "linuxremoteplayer_${VERSION}_all.deb" "$DEB_URL"

echo "$SHA256 linuxremoteplayer_${VERSION}_all.deb" | sha256sum -c -

apt-get install -y --reinstall "./linuxremoteplayer_${VERSION}_all.deb"

# Restart services depending on mode
if systemctl is-enabled linuxremoteplayer.service &>/dev/null; then
    systemctl restart linuxremoteplayer.service || true
else
    # find user services
    for d in /run/user/*; do
        if [ -d "$d" ]; then
            uid=$(basename "$d")
            sudo -u "#$uid" XDG_RUNTIME_DIR="$d" systemctl --user restart linuxremoteplayer.service || true
        fi
    done
fi

# Update uBlock Origin Lite
echo "Actualizando uBlock Origin Lite..."
mkdir -p /opt/linuxremoteplayer/extensions
curl -fsSL "https://github.com/uBlockOrigin/uBOL-home/releases/latest/download/uBOLite_mv3.zip" -o /tmp/ubol.zip
unzip -qo /tmp/ubol.zip -d /opt/linuxremoteplayer/extensions/ubol || true
rm -f /tmp/ubol.zip
# Ownership is fixed by install.sh or by the running service.
# Let's fix it here for good measure
if [ -n "$SUDO_USER" ]; then
    chown -R "$SUDO_USER":"$SUDO_USER" /opt/linuxremoteplayer/extensions
fi
INNEREOF
chmod 755 /usr/local/bin/lrp-update

echo "================================================="
echo " Instalación base completada. "
echo " Ejecuta: sudo lrp-setup para configurar tu modo."
echo "================================================="
